# logrotate/logger.service
[Unit]
Description=My logger shell script

[Service]
ExeStart=/usr/local/bin/logger.sh

[Install]
wantedBy=multi-user.target
